# api-gateway/docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
autostart=true
autorestart=true
user=root
redirect_stderr=true

[program:default-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --sleep=3 --tries=3 --timeout=90
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/worker.log
stderr_logfile=/var/log/worker.err.log
stopwaitsecs=3600

[program:messages-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work rabbitmq_message_notify --queue=message_notify_queue
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/message_notify_queue.log
stderr_logfile=/var/log/message_notify_queue.err.log
stopwaitsecs=3600

[program:reverb]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan reverb:start --debug
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/reverb.log
stderr_logfile=/var/log/reverb.err.log
stopwaitsecs=3600

[program:laravel-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan schedule:run
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/laravel-schedule.log
stderr_logfile=/var/log/laravel-schedule.err.log
stopwaitsecs=3600